# infra/bpf-runner/assay-update.service
[Unit]
Description=Update assay CLI to the latest GitHub release
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/update-assay-latest
